# PSPi/Configs/PSPi.service
[Unit]
Description=Controller Service for PSPi

[Service]
Type=oneshot
ExecStart=bash /home/pi/PSPi/boot.sh
ExecStop=bash /home/pi/PSPi/shutdown.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
